# .config/sv/emacs/run
#!/bin/sh

cd
exec emacs --fg-daemon &>> "$XDG_DATA_HOME/sv/generic.log"
