# .config/sv/zsh/finish
#!/bin/sh

kill -HUP $_RUNSVDIR_USER_PID
