# .config/sv/slstatus/run
#!/bin/sh

exec slstatus &>> "$XDG_DATA_HOME/sv/generic.log"
