# .config/service/offlineimap/run
#!/bin/sh

while true; do
snooze -m *
echo "LOG" >> ~/.tmpfile
done
