# .config/service/offlineimap/run
#!/bin/sh

while true; do
snooze -S0 -M0 -H '/2'

offlineimap
done
