# .config/sv/pipewire/run
#!/bin/sh

exec pipewire &>> "$XDG_DATA_HOME/sv/generic.log"
