# .config/sv/i3scratchpad_musicplayer/run
#!/bin/sh

cd
exec kitty --name i3scratchpad_musicplayer rmpc

