# .config/sv/i3scratchpad_musicplayer/run
#!/bin/sh

cd
exec xterm -name i3scratchpad_musicplayer rmpc

