# .config/sv/i3scratchpad_musicplayer/run
#!/bin/sh

exec kitty --name i3scratchpad_musicplayer rmpc

