# .config/sv/i3scratchpad_musicplayer/run
#!/bin/sh

cd
exec kitty --name i3scratchpad_musicplayer rmpc &>> "$XDG_DATA_HOME/sv/generic.log"

