# .config/sv/i3/finish
#!/bin/sh

kill -HUP $_RUNSVDIR_X11_PID
