# .config/sv/sxhkd/run
#!/bin/sh

cd "$HOME"
exec sxhkd &>> "$XDG_DATA_HOME/sv/generic.log"
