# .config/sv/wallpaper_shuffler/run
#!/bin/sh

exec "$HOME/.local/libexec/wallpaper_shuffler" &>> "$XDG_DATA_HOME/sv/generic.log"
