# .config/sv/mailsync/run
#!/bin/sh

exec "$HOME/.local/libexec/mailsync" &>> "$XDG_DATA_HOME/sv/generic.log"
