# .config/sv/mailsync/run
#!/bin/sh

while true; do
	snooze -H* -M/10
	offlineimap -l "$XDG_STATE_HOME/offlineimap.log"
done
