# .config/sv/i3scratchpad_dropdown/run
#!/bin/sh

cd
exec st -n i3scratchpad_dropdown

