# .config/sv/i3scratchpad_dropdown/run
#!/bin/sh

cd
NOHEADER=t DROPDOWN=t exec st -n i3scratchpad_dropdown &>> "$XDG_DATA_HOME/sv/generic.log"

