# .config/sv/i3scratchpad_dropdown/run
#!/bin/sh

exec st -n i3scratchpad_dropdown

