# .config/sv/i3scratchpad_dropdown/run
#!/bin/sh

cd
NOHEADER=t DROPDOWN=t exec st -n scratchpad_dropdown -e tmux -L dropdown >> "$XDG_DATA_HOME/sv/generic.log" 2>&1

