# .config/sv/i3scratchpad_dropdown/run
#!/bin/sh

cd
NOHEADER=t SMALLPROMPT=t exec xterm -name i3scratchpad_dropdown

