# .config/sv/i3scratchpad_dropdown/run
#!/bin/sh

cd
NOHEADER=t SMALLPROMPT=t exec st -n i3scratchpad_dropdown -f "Cozette:pixelsize=32:antialias=false:autohint=true"

