# .config/sv/picom/run
#!/bin/sh

exec picom  >> "$XDG_DATA_HOME/sv/generic.log" 2>&1
