# .config/sv/xbanish/run
#!/bin/sh

exec xbanish -t 5 &>> "$XDG_DATA_HOME/sv/generic.log"

