# .config/sv/mpd/run
#!/bin/sh

exec mpd --no-daemon  >> "$XDG_DATA_HOME/sv/generic.log" 2>&1
