# .config/sv/syndaemon/run
#!/bin/sh

exec syndaemon -i 0.5 -R -t -k
