# .config/sv/i3scratchpad_btop/run
#!/bin/sh

st -n i3scratchpad_btop -e btop
