# .config/sv/i3scratchpad_btop/run
#!/bin/sh

cd
exec st -n i3scratchpad_btop -e btop &>> "$XDG_DATA_HOME/sv/generic.log"
