# .config/sv/i3scratchpad_btop/run
#!/bin/sh

cd
st -n i3scratchpad_btop -e btop
