# .config/sv/i3scratchpad_btop/run
#!/bin/sh

cd
exec st -n i3scratchpad_btop -f "Cozette:pixelsize=24:antialias=false:autohint=true" -e btop
