# .config/sv/scratchpad_dropdown/run
#!/bin/sh

cd
NOHEADER=t DROPDOWN=t exec st -n scratchpad_dropdown -e tmux -L dropdown new -s drop >> "$XDG_DATA_HOME/sv/generic.log" 2>&1

