# .config/sv/acpi_monitor/run
#!/bin/sh

exec "$HOME/.local/libexec/acpi_monitor"
